# lumakit.service.example
[Unit]
Description=LumaKit AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/LumaKit
EnvironmentFile=/path/to/LumaKit/.env
ExecStart=/usr/bin/python3 -m lumakit serve
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
